# .runit/sv/picom/run
#!/bin/sh
# shellcheck disable=1091

exec 2>&1

[ -r ./conf ] && . ./conf

# shellcheck disable=2086
# shfmt:ignore
exec picom $OPTS
