# .runit/sv/runsvdir-jaagr/run
#!/bin/sh

if [ ! -d /run/runit.jaagr ]; then
  mkdir -p /run/runit.jaagr
  chown jaagr:jaagr /run/runit.jaagr
fi

export HOME=/home/jaagr
groups=$(id -Gn jaagr | tr ' ' ':')

exec 2>&1
exec chpst -u "jaagr:$groups" sh -l -c "exec runsvdir /home/jaagr/.local/sv 'log: ...........................................................................................................................................................................................................................................................................................................................................................................................................'" jaagr
