# .runit/sv/sxhkd/finish
#!/bin/sh

"$XDG_CONFIG_HOME"/sxhkd/eventhandler stop -f
