# .runit/sv/dbus-session/run
#!/bin/sh

exec dbus-daemon --session --nofork --nopidfile 2>&1
