# .runit/sv/bspwm/finish
#!/bin/sh

"$XDG_CONFIG_HOME"/bspwm/eventhandler stop -f
